# eos-license-service.socket
[Unit]
Description=EndlessOS License Service Socket
Wants=network.target

[Socket]
ListenStream=3010

[Install]
WantedBy=graphical.target
