# eos-license-service.service.in
[Unit]
Description=EndlessOS License Service
Wants=network.target

[Service]
Environment=NODE_ENV=production
ExecStart=/usr/bin/nodejs @pkgdatadir@/licenseApp.js
