# eos-license-service.service.in
[Unit]
Description=EndlessOS License Service
Wants=network.target

[Service]
ExecStart=/usr/bin/gjs -I @pkgdatadir@ @pkgdatadir@/licenseApp.js
